# orderly-bot.service
[Unit]
Description=Orderly Discord Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/harold/.openclaw/workspace/projects/orderly_docs
Environment="PYTHONPATH=/home/harold/.openclaw/workspace/projects/orderly_docs/src"
Environment="PATH=/home/harold/.local/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/harold/.openclaw/workspace/projects/orderly_docs/venv/bin/python /home/harold/.openclaw/workspace/projects/orderly_docs/src/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
